# src/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --issue-file=/etc/issue --autologon root %I $TERM
